# systemd/ap-disable.timer
[Unit]
Description=Disable travel router AP at night

[Timer]
OnCalendar=*-*-* 02:00:00
Persistent=true

[Install]
WantedBy=timers.target
