# systemd/ap-enable.timer
[Unit]
Description=Re-enable travel router AP in the morning

[Timer]
OnCalendar=*-*-* 07:00:00
Persistent=true

[Install]
WantedBy=timers.target
