# systemd/ap-enable.service
[Unit]
Description=Re-enable travel router AP
After=NetworkManager.service network.target
ConditionPathExists=/usr/local/bin/ap-schedule.sh

[Service]
Type=oneshot
ExecStart=/usr/local/bin/ap-schedule.sh enable
EnvironmentFile=-/etc/default/travel-router
