# systemd/daily-digest.timer
[Unit]
Description=Daily router digest at 08:00

[Timer]
OnCalendar=*-*-* 08:00:00
Persistent=true
RandomizedDelaySec=15min

[Install]
WantedBy=timers.target
