# systemd/generate-bandwidth-report.service
[Unit]
Description=Generate HTML bandwidth report from vnStat
Wants=vnstat.service
After=vnstat.service

[Service]
Type=oneshot
ExecStart=/usr/local/bin/generate-bandwidth-report.sh
StandardOutput=journal
StandardError=journal
