# systemd/split-tunnel.service
[Unit]
Description=Domain-based split tunnel routing
After=network-online.target tailscaled.service
Wants=tailscaled.service

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/local/bin/apply-split-tunnel.sh

[Install]
WantedBy=multi-user.target
