# systemd/tune-cake.timer
[Unit]
Description=Weekly CAKE bandwidth auto-tune

[Timer]
OnCalendar=Sun 04:00
Persistent=true
RandomizedDelaySec=15min

[Install]
WantedBy=timers.target
