# systemd/update-router.service
[Unit]
Description=Auto-update travel router scripts from GitHub
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/update-router.sh
StandardOutput=journal
StandardError=journal
