# systemd/ups-monitor.timer
[Unit]
Description=UPS battery check every 5 minutes

[Timer]
OnBootSec=60
OnUnitActiveSec=5min

[Install]
WantedBy=timers.target
