# systemd/wan-watchdog.service
[Unit]
Description=WAN connectivity watchdog with graduated recovery
After=network.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/wan-watchdog.sh
Restart=on-failure
RestartSec=10
