# firstboot/firstboot.service
[Unit]
Description=Travel Router first-boot setup wizard
# imager-compat.service must complete before we start: it neutralises
# firstrun.sh and writes the SSH key from Raspberry Pi Imager.  The
# Before= in imager-compat.service is not enough on its own when
# firstboot.service lives in multi-user.target — make the ordering
# explicit here too so systemd enforces it regardless of target changes.
After=network.target network-online.target avahi-daemon.service imager-compat.service
Wants=network-online.target
RequiresMountsFor=/opt
ConditionPathExists=!/var/lib/travel-router/firstboot-done

[Service]
Type=simple
ExecStartPre=/bin/mkdir -p /var/lib/travel-router
ExecStart=/usr/bin/python3 /opt/pi-travel-router/firstboot/server.py
Restart=on-failure
RestartSec=5
TimeoutStartSec=infinity
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
